# s6-overlay/s6-rc.d/sleev/finish
#!/command/with-contenv sh
if [ ! -x "/home/opencode/.local/share/sleev/gateway/current" ]; then
  exit 125
fi
exit "$1"
