# s6-overlay/s6-rc.d/hermes/run
#!/command/with-contenv sh
exec s6-setuidgid opencode env \
  API_SERVER_ENABLED=true \
  API_SERVER_HOST="${API_SERVER_HOST:-0.0.0.0}" \
  API_SERVER_PORT="${HERMES_PORT:-8642}" \
  HERMES_HOME="${HERMES_HOME:-/home/opencode/.hermes}" \
  hermes gateway run
